# deploy/celery.service
[Unit]
Description=ERP Nexus - Celery Worker
After=network.target redis-server.service

[Service]
Type=forking
User=erp
Group=www-data
WorkingDirectory=/opt/erpnexus/erp_backend

# Ambiente
EnvironmentFile=/opt/erpnexus/erp_backend/.env
Environment="PATH=/opt/erpnexus/erp_backend/.venv/bin"

# Comando
ExecStart=/opt/erpnexus/erp_backend/.venv/bin/celery \
    -A config \
    worker \
    --loglevel=info \
    --logfile=/opt/erpnexus/logs/celery_worker.log \
    --pidfile=/run/celery_worker.pid \
    --concurrency=4

# Reiniciar policies
Restart=on-failure
RestartSec=10s

# Limites
StandardOutput=journal
StandardError=journal
SyslogIdentifier=celery

[Install]
WantedBy=multi-user.target
